# launchd/com.buildguard.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.buildguard.frontend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/tangkai/Desktop/BuildGuard/frontend/node_modules/vite/bin/vite.js</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>5173</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/tangkai/Desktop/BuildGuard/frontend</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>5</integer>
  <key>StandardOutPath</key>
  <string>/Users/tangkai/Desktop/BuildGuard/.runtime/logs/frontend-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tangkai/Desktop/BuildGuard/.runtime/logs/frontend-launchd-error.log</string>
</dict>
</plist>
